# vms/efi/mkosi/ubuntu/x86_64/mkosi.extra/etc/udev/rules.d/90-etcd-network-tuning.rules
ACTION=="add", SUBSYSTEM=="net", SUBSYSTEMS=="pci|xen|vmbus" RUN+="/usr/local/sbin/etcd-network-tuning.sh $name
