# sysext/kubernetes/kubeadm/confext-overlay/etc/systemd/system/containerd.service.d/memory-pressure.conf
[Service]
# https://www.freedesktop.org/software/systemd/man/systemd.exec.html#OOMScoreAdjust=
OOMScoreAdjust=-999
