# sysext/kubernetes/kubeadm/confext-overlay/etc/profile.d/kubernetes.sh
export KUBECONFIG=/etc/kubernetes/admin.conf
export PATH="/opt/kubernetes/bin:/opt/bin:$PATH"
